# beifen/beifen.service
[Unit]
Description=自定义备份脚本

[Service]
ExecStart=/bin/bash /home/user/scripts/mdweb.sh
